# blob-wrangler.openrc
#!/sbin/openrc-run

description="Android vendor firmware extraction"

depend() {
  after mount
  after tmpfs
  need localmount
}

start_pre() {
  echo "Extracting binary firmware..."
}

start() {
  if [ ! -f /var/lib/blob-wrangler/status.json ]; then
    /usr/bin/blob-wrangler
    /sbin/reboot
  fi
}

stop() {
  return 0
}
