# droid-juicer.service
[Unit]
Description=Android vendor firmware extraction
Wants=sysinit.target
After=-.mount tmp.mount boot.mount
Before=plymouth-quit.service display-manager.service

[Service]
Type=oneshot
Environment=RUST_LOG=debug
ExecStartPre=-/usr/bin/plymouth display-message --text="Extracting binary firmware..."
ExecStart=/usr/bin/droid-juicer

[Install]
WantedBy=basic.target
